# al-nassr-news-monitor.service
[Unit]
Description=AlNassr News Monitor Service
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/AlNassr-News-Monitor
ExecStart=/opt/AlNassr-News-Monitor/venv/bin/python /opt/AlNassr-News-Monitor/src/main.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/opt/AlNassr-News-Monitor/logs/service.log
StandardError=append:/opt/AlNassr-News-Monitor/logs/service.log
Environment="PATH=/opt/AlNassr-News-Monitor/venv/bin"
EnvironmentFile=/opt/AlNassr-News-Monitor/.env

[Install]
WantedBy=multi-user.target
